# ansible-bdas/roles/anaconda/files/etc/profile.d/anaconda.sh
export PATH=/opt/anaconda/bin:/opt/anaconda/lib:$PATH